# overlay/etc/cron.d/wordpress-cron
#!/bin/bash
*/5 * * * * www-data /usr/bin/curl -A 'Mozilla/5.0' 'http://127.0.0.1/wp-cron.php' >/dev/null 2>&1
